# backups/systemd/backup-heartbeat.timer
[Unit]
Description=Backup Systems Daily Heartbeat Timer
Requires=backup-heartbeat.service

[Timer]
# Run daily at 3:30 AM (30 minutes after backup completion)
OnCalendar=*-*-* 03:30:00
Persistent=true

# Allow some jitter to avoid conflicts
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
