# backups/systemd/postgres-diff-backup.timer
[Unit]
Description=PostgreSQL Differential Backup Timer (Daily at 3 AM)
Documentation=https://pgbackrest.org/
Requires=postgres-diff-backup.service

[Timer]
# Run daily at 3:00 AM (Monday through Saturday; Sunday runs full backup instead)
OnCalendar=Mon,Tue,Wed,Thu,Fri,Sat *-*-* 03:00:00

# If system was off, run backup on next boot
Persistent=true

# Randomize by 5 minutes to avoid conflicts
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
