# backups/systemd/postgres-full-backup.service
[Unit]
Description=PostgreSQL Full Backup with pgBackRest
Documentation=https://pgbackrest.org/
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/keith/digitalocean-stack/db

# Run FULL backup (weekly)
ExecStart=/usr/bin/docker exec db /usr/local/bin/pgbackrest-wrapper --stanza=main --type=full backup

# On failure, send Discord alert
ExecStartPost=/bin/bash -c 'if [ $EXIT_STATUS -ne 0 ]; then /home/keith/digitalocean-stack/backups/monitoring/send-failure-alert.sh "Full Backup" "$EXIT_STATUS"; fi'

# On success, log to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=postgres-full-backup

# Restart policy
Restart=on-failure
RestartSec=5m

# Resource limits (full backups can use more resources)
Nice=10
IOSchedulingClass=best-effort
CPUQuota=75%

# Security
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
