# backups/monitoring/heartbeat/backup-heartbeat.timer
[Unit]
Description=Backup Systems Daily Heartbeat Timer
Requires=backup-heartbeat.service

[Timer]
# Run daily at 5:00 AM (2 hours after backup start, ensuring completion)
# Timeouts: postgres-diff=1h, postgres-full=2h, foundry=1.5h
OnCalendar=*-*-* 05:00:00
Persistent=true

# Allow some jitter to avoid conflicts
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
