# backups/foundry/foundry-backup.service
[Unit]
Description=Foundry VTT Data Backup
After=network-online.target
Wants=network-online.target
OnFailure=backup-failure-alert@foundry-backup.service

[Service]
Type=oneshot
# Run as stack user to access Foundry data and restic cache without permission issues
User=@STACK_USER@
Group=@STACK_USER@
WorkingDirectory=@STACK_DIR@/backups/foundry

# Run backup via wrapper script
ExecStart=@STACK_DIR@/backups/foundry/backup-wrapper.sh

# Timeout after 1.5 hours (restic backups can take time)
TimeoutStartSec=5400

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foundry-backup

# Restart policy
Restart=no

# Resource limits
CPUQuota=50%
MemoryMax=2G

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@STACK_DIR@/backups @STACK_DIR@/foundry/data /home/@STACK_USER@/.cache/restic
