# backups/foundry/foundry-backup.service
[Unit]
Description=Foundry VTT Data Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=@STACK_USER@
Group=@STACK_USER@
WorkingDirectory=@STACK_DIR@/backups/foundry
ExecStart=@STACK_DIR@/backups/foundry/backup-wrapper.sh
ExecStartPost=/bin/bash -c 'if [ $EXIT_STATUS -ne 0 ]; then @STACK_DIR@/backups/monitoring/send-failure-alert.sh "Foundry Backup" $EXIT_STATUS; fi'

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@STACK_DIR@/backups/logs @STACK_DIR@/foundry/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foundry-backup

# Resource limits
CPUQuota=50%
MemoryMax=2G

# Failure handling
Restart=no
